# Dockerfile/zabbix-3.0/container-files-base/etc/supervisor.d/php-fpm.conf
[program:php-fpm]
command = /opt/remi/php70/root/usr/sbin/php-fpm
autorestart = true
